# systemd/solomon-claw.service
[Unit]
Description=SolomonClaw AI Telegram Bot
After=network.target

[Service]
Type=simple
User=solomon-claw
WorkingDirectory=/opt/solomon-claw
EnvironmentFile=/opt/solomon-claw/.env
ExecStart=/usr/bin/node /opt/solomon-claw/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/opt/solomon-claw/data
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
